# twitter-score-updater/docker/crontab
MAILTO=root

45 16 * * * root java -jar /score-updater-assembly-0.0.1.jar >> /var/log/updater.log 2>&1
